# deploy/portfolio.service
[Unit]
Description=Portfolio FastAPI Application
After=network.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/var/www/portfolio
EnvironmentFile=/var/www/portfolio/.env
ExecStart=/var/www/portfolio/venv/bin/gunicorn backend.main:app \
    --workers 2 \
    --worker-class uvicorn.workers.UvicornWorker \
    --bind 127.0.0.1:8000 \
    --timeout 120 \
    --access-logfile /var/log/portfolio/access.log \
    --error-logfile /var/log/portfolio/error.log
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
